# soft-hwclock-tick.service.in
[Unit]
Description=Save hwclock by timer

[Service]
Type=oneshot
ExecStart=@BINDIR@/soft-hwclock tick > /dev/null

